# ansible/roles/master/files/etc/cron.d/backup_scratch_cm
@weekly sbatch /home/backup/bin/backup_scratch_cm_slurm.sh
